# profile.d/99-paths.sh
export PATH="$(printf "%s:" "${paths[@]}")${PATH}"
unset paths
